# systemd/zsys-commit.service
[Unit]
Description=Mark current ZSYS boot as successful
After=default.target zsysd.socket zsysd.service
Requires=default.target zsysd.socket zsysd.service

# We can't run it in a container
ConditionVirtualization=!container

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/zsysctl boot commit
ExecStop=/sbin/zsysctl boot update-lastused

# Some daemon restrictions
NoNewPrivileges=true
PrivateNetwork=true

[Install]
WantedBy=default.target