# systemd/zsysd.socket
[Unit]
Description=Socker activation for zsys daemon
PartOf=zsyd.service

[Socket]
ListenStream=/run/zsysd.sock

[Install]
WantedBy=sockets.target